# resources/supervisord.conf
[supervisord]
nodaemon=true

[program:vpnserver]
command='/usr/local/vpnserver/runner.sh'